# deploy/com.tenda.monitor.poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.tenda.monitor.poller</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/bin/node</string>
      <string>/Users/michaeltawfik/Documents/tenda_traffic/packages/poller/dist/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/michaeltawfik/Documents/tenda_traffic</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/michaeltawfik/Library/Logs/tenda-monitor/poller.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/michaeltawfik/Library/Logs/tenda-monitor/poller.error.log</string>
  </dict>
</plist>
